# contrib/systemd/dovecot-fts-optimize.timer.in
[Unit]
Description=Optimize Dovecot FTS Index

[Timer]
OnCalendar=@@index_frequency@@
Persistent=true

[Install]
WantedBy=timers.target
